# benchmarks/nginx-zmq/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=nginx -g 'daemon off;'

[program:backend]
command=./host uppercase.xclbin
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0